# plem.sapphic.site/etc/systemd/system/secure-time-sync.timer
[Unit]
Description=Secure time sync 1 hour
Requires=secure-time-sync.service

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target