# plem.sapphic.site/etc/systemd/system/gitea.service.d/override.conf
[Unit]
StartLimitInterval=0

[Service]
ReadWritePaths=/mnt/pleroma_storage/gitea
ReadOnlyPaths=/etc/gitea
MemoryHigh=475M
MemoryMax=520M

BindPaths=/mnt/pleroma_storage/gitea