# plem.sapphic.site/etc/systemd/system/chronyd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/chronyd -F 2 $OPTIONS